# etc/systemd/lfsd-sync.service
[Unit]
Description=LFSd Sync Receitas

[Service]
Type=oneshot
Environment=LFSD_CHANNEL=stable
ExecStart=/usr/bin/lfsd sync
